# emulab-devel-new/emulab-devel/clientside/tmcc/linux/docker/dockerfiles/alpine/fs/etc/service/rsyslog/run
#!/bin/sh

exec /usr/sbin/rsyslogd -n $RSYSLOGD_OPTIONS
