# emulab-devel-new/emulab-devel/clientside/tmcc/linux/docker/dockerfiles/common/fs/etc/service/dockerentrypoint/finish
#!/bin/sh

#
# We never restart this service.  If it stops, presumably user wants to
# login to see what happened.
#

sv down dockerentrypoint
