# emulab-devel-new/emulab-devel/install/oports/net-snmp4/files/snmpd.in
#!/bin/sh
#
# $FreeBSD: ports/net-mgmt/net-snmp4/files/snmpd.in,v 1.4 2010/03/27 00:14:24 dougb Exp $
#
# PROVIDE: snmpd
# REQUIRE: DAEMON
#
# Add the following lines to /etc/rc.conf to enable this service:
#
# snmpd_enable (bool):	Set to NO by default. Set it to YES to enable snmpd.
#

. /etc/rc.subr

name=snmpd
rcvar=${name}_enable

command=/usr/local/sbin/${name}
pidfile=/var/run/${name}.pid

load_rc_config ${name}

: ${snmpd_enable="NO"}
: ${snmpd_flags="-P /var/run/${name}.pid"}

run_rc_command "$1"
