# emulab-devel-new/emulab-devel/install/ports/emulab-pubsub/pkg-plist
libexec/clusterd
libexec/pubsubd
libexec/pubsub_forward
include/pubsub/elvin_hash.h
include/pubsub/pubsub.h
lib/libpubsub.a
lib/libpubsub.la
lib/libpubsub.so
lib/libpubsub.so.1
lib/libpubsub.so.1.0.0
lib/libpubsub_r.a
lib/libpubsub_r.la
lib/libpubsub_r.so
lib/libpubsub_r.so.1
lib/libpubsub_r.so.1.0.0
etc/rc.d/2.pubsubd.sh
