# emulab-devel-new/emulab-devel/clientside/tmcc/linux/docker/dockerfiles/debian8/fs/etc/service/testbed/finish
#!/bin/sh

exec /etc/init.d/testbed stop

exit 1
