# emulab-devel-new/emulab-devel/clientside/tmcc/linux/docker/dockerfiles/alpine/fs/etc/service/pubsubd/finish
#!/bin/sh

if [ -f /etc/init.d/pubsubd ]; then
    exec /etc/init.d/pubsubd stop
else
    sv down pubsubd
fi

exit 0
