# emulab-devel-new/emulab-devel/clientside/tmcc/linux/docker/dockerfiles/centos/fs/etc/service/sshd/run
#!/bin/sh

# Never run if systemd is running; it will do this.
[ -f /bin/systemctl ] && /bin/systemctl show-environment >/dev/null 2>&1 && exit 0

mkdir -p /var/run/sshd

exec /usr/sbin/sshd -D
