# emulab-devel-new/emulab-devel/clientside/tmcc/linux/docker/dockerfiles/alpine/fs/etc/service/testbed/finish
#!/bin/sh

if [ -f /usr/local/etc/emulab/initscripts/testbed ]; then
    exec /usr/local/etc/emulab/initscripts/testbed stop
fi

exit 0
