# emulab-devel/clientside/tmcc/linux/docker/dockerfiles/debian/fs/etc/service/getty/run
#!/bin/sh

exec getty 115200 console linux
