# emulab-devel/clientside/tmcc/linux/docker/dockerfiles/alpine/fs/etc/service/tbprepare/finish
#!/bin/sh

if [ -f /usr/local/etc/emulab/initscripts/tbprepare ]; then
    exec /usr/local/etc/emulab/initscripts/tbprepare stop
fi

exit 0
