# emulab-devel/clientside/tmcc/linux/openvz/prunelvms.cron
#! /bin/bash

#
# Prune our lvms.
#
/usr/local/etc/emulab/prunelvms -d >/var/emulab/logs/prunelvms 2>&1
exit 0
