# emulab-devel-new/emulab-devel/clientside/tmcc/linux/xen/udev/40-bridge-network-interface.rules
# do not do this for ifb devices
ACTION=="add", SUBSYSTEM=="net", KERNEL!="ifb*", RUN+="bridge-network-interface"
