# emulab-devel/clientside/tmcc/linux/docker/dockerfiles/alpine/fs/etc/service/getty/run
#!/bin/sh

exec agetty 115200 console linux
