# emulab-devel-new/emulab-devel/clientside/tmcc/ubuntu20/01-disable-bluefield-nvme.rules
ACTION=="add", SUBSYSTEM=="pci", ATTR{vendor}=="0x15b3", ATTR{device}=="0x6001", RUN+="/bin/sh -c 'echo 1 > /sys/$env{DEVPATH}/remove'"
